# systemd/remindd.service
[Unit]
Description=remindme Daemon

[Service]
Type=simple
ExecStart=/usr/local/bin/remindd

[Install]
WantedBy=multi-user.target
